# data/systemd-user/snapd.aa-prompt-ui.service.in
[Unit]
Description=snapd UI prompt
After=gnome-session-initialized.target

[Service]
Type=simple
ExecStart=@libexecdir@/snapd/snapd-aa-prompt-ui
