# data/systemd-user/snapd.session-agent.socket
[Unit]
Description=REST API socket for snapd user session agent

[Socket]
ListenStream=%t/snap-session.socket

[Install]
WantedBy=sockets.target
