# data/systemd-user/snapd.aa-prompt-ui.service.in
[Unit]
Description=snapd UI prompt
ConditionPathExists=/sys/kernel/security/apparmor/.notify

[Service]
Type=simple
ExecStart=@libexecdir@/snapd/snapd-aa-prompt-ui
