# data/systemd-user/snapd.aa-prompt-ui.service.in
[Unit]
Description=snapd UI prompt

[Service]
Type=dbus
BusName=io.snapcraft.Prompt
ExecStart=@bindir@/snapd-aa-prompt-ui
