# data/systemd-user/snapd.session-agent.service.in
[Unit]
Description=snapd user session agent
Requires=snapd.session-agent.socket

[Service]
Type=notify
ExecStart=@bindir@/snap userd --agent
