# core-initrd/26.04/factory/usr/lib/systemd/system/plymouth-start.service.d/core-override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/plymouthd --mode=boot --pid-file=/run/plymouth/pid --attach-to-session --ignore-serial-consoles
