# data/systemd/snapd.check-chooser.service.in
[Unit]
Description=Check whether the chooser menu should be invoked
Before=snapd.service
# don't run on classic or uc16/uc18
ConditionKernelCommandLine=snapd_recovery_mode

[Service]
Type=simple
ExecStart=@libexecdir@/snapd/snap-bootstrap check-chooser
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
