# data/systemd/snap.mount.in
# ensure @SNAP_MOUNT_DIR@ MS_SHARED so that lxd works correctly (LP: #1668659)
[Unit]
Description=Setup @SNAP_MOUNT_DIR@ mountpoint

[Mount]
What=@SNAP_MOUNT_DIR@
Where=@SNAP_MOUNT_DIR@
Type=none
Options=bind,shared

[Install]
WantedBy=multi-user.target