# data/systemd/snapd.system-services-pre.target
[Unit]
Description=Preparation for snap system services
Wants=snapd.apparmor.service
Wants=network.target
