# data/systemd/snapd.autoimport.service.in
[Unit]
Description=Auto import assertions from block devices
After=snapd.service snapd.socket snapd.seeded.service
# don't run on classic
ConditionKernelCommandLine=|snap_core
# TODO:UC20: only enable this in run mode?
ConditionKernelCommandLine=|snapd_recovery_mode

[Service]
Type=oneshot
ExecStart=@bindir@/snap auto-import
EnvironmentFile=-/var/lib/snapd/environment/snapd.conf

[Install]
WantedBy=multi-user.target
# This cannot be started on firstboot
# X-Snapd-Snap: do-not-start
