# data/systemd/snapd.system-shutdown.service.in
[Unit]
Description=Ubuntu core (all-snaps) system shutdown helper setup service
Before=umount.target
DefaultDependencies=false
# don't run on classic
ConditionKernelCommandLine=snap_core
# don't run if system-shutdown isn't there
ConditionPathExists=@libexecdir@/snapd/system-shutdown
# X-Snapd-Snap: do-not-start

[Service]
Type=oneshot
ExecStart=/bin/sh -euc 'mount /run -o remount,exec; mkdir -p /run/initramfs; cp @libexecdir@/snapd/system-shutdown /run/initramfs/shutdown'

[Install]
WantedBy=final.target
