# data/systemd/snapd.seeded.service.in
[Unit]
Description=Wait until snapd is fully seeded
After=snapd.service snapd.socket
Requires=snapd.socket

[Service]
Type=oneshot
ExecStart=@bindir@/snap wait system seed.loaded
RemainAfterExit=true

[Install]
RequiredBy=snapd.service
# This will be started when snapd.service is started so no need
# for the snapd snap to do anything about it.
# X-Snapd-Snap: do-not-start
