# data/systemd/snapd.seeded.service.in
[Unit]
Description=Wait until snapd is fully seeded
After=snapd.socket
Requires=snapd.socket

[Service]
Type=oneshot
ExecStart=@bindir@/snap wait system seed.loaded
RemainAfterExit=true

[Install]
WantedBy=multi-user.target cloud-final.service
# This is handled special in snapd
# X-Snapd-Snap: do-not-start
