# data/systemd/snapd.mounts.snapd.target
[Unit]
Description=Mounted snapd snaps
Wants=snapd.mounts-pre.target
# Started from snapd-apparmor.service
# X-Snapd-Snap: do-not-start
