# data/systemd/snapd.aa-prompt-listener.service.in
[Unit]
Description=Userspace listener for prompt events
After=snapd.socket
ConditionPathExists=/sys/kernel/security/apparmor/.notify

[Service]
ExecStart=@libexecdir@/snapd/snapd-aa-prompt-listener
EnvironmentFile=-@SNAPD_ENVIRONMENT_FILE@
Restart=on-failure
Type=simple

[Install]
WantedBy=multi-user.target
