# data/systemd/snapd.recovery-chooser-trigger.service.in
[Unit]
Description=Wait for the Ubuntu Core chooser trigger
Wants=getty-pre.target
Before=getty-pre.target
# don't run on classic or uc16/uc18
ConditionKernelCommandLine=snapd_recovery_mode
# only run when there are input devices
ConditionPathExistsGlob=/dev/input/event*

[Service]
# blocks the service startup until a trigger is detected or a timeout is hit
Type=oneshot
ExecStart=@libexecdir@/snapd/snap-bootstrap recovery-chooser-trigger
RemainAfterExit=true

[Install]
WantedBy=multi-user.target

# started on boot only
# X-Snapd-Snap: do-not-start
