# data/systemd/snapd.system-shutdown.service.in
[Unit]
Description=Ubuntu core (all-snaps) system shutdown helper setup service
Before=umount.target
DefaultDependencies=false
# don't run on classic
ConditionKernelCommandLine=|snap_core
ConditionKernelCommandLine=|snapd_recovery_mode
# don't run if system-shutdown isn't there
ConditionPathExists=@libexecdir@/snapd/system-shutdown
# X-Snapd-Snap: do-not-start

[Service]
Type=oneshot
ExecStart=/bin/mount /run -o remount,exec
ExecStart=/bin/mkdir -p /run/initramfs
ExecStart=/bin/cp @libexecdir@/snapd/system-shutdown /run/initramfs/shutdown

[Install]
WantedBy=final.target
