# data/systemd/snapd.failure.service.in
[Unit]
Description=Failure handling of the snapd snap

[Service]
Type=oneshot
ExecStart=@libexecdir@/snapd/snap-failure run
RemainAfterExit=true
