# data/systemd/snapd.gpio-chardev-setup.target
[Unit]
Description=Synchronization target for gpio-chardev interface setup
# Started from gpio-chardev setup services
# X-Snapd-Snap: do-not-start
