# data/systemd/snap-repair.service.in
[Unit]
Description=Automatically fetch and run repair assertions
ConditionPathExistsGlob=@SNAP_MOUNT_DIR@/*/current
Documentation=man:snap(1)

[Service]
Type=oneshot
ExecStart=@libexecdir@/snapd/snap-repair run
Environment=SNAP_REAPIR_FROM_TIMER=1
