# data/systemd/snapd.seed-done.service.in
[Unit]
Description=Wait until snapd is fully seeded
After=snapd.service snapd.socket

[Service]
Type=oneshot
ExecStart=@bindir@/snap wait system seed.done
RemainAfterExit=true

[Install]
WantedBy=multi-user.target cloud-final.service

