# data/systemd/snapd.service.in
[Unit]
Description=Snappy daemon
Requires=snapd.socket
StartLimitBurst=10
StartLimitIntervalSec=105s

[Service]
Nice=-5
OOMScoreAdjust=-900
ExecStart=@libexecdir@/snapd/snapd
EnvironmentFile=-@SNAPD_ENVIRONMENT_FILE@
Restart=always
Type=notify
TimeoutStartSec=10s

[Install]
WantedBy=multi-user.target
