# data/systemd/snapd.autoimport.service.in
[Unit]
Description=Auto import assertions from block devices
After=snapd.service snapd.socket
# don't run on classic
ConditionKernelCommandLine=snap_core

[Service]
Type=oneshot
ExecStart=@bindir@/snap auto-import

[Install]
WantedBy=multi-user.target
