# data/systemd/snap-userd.service.in
[Unit]
Description=Snap User Session Service

[Service]
ExecStart=@bindir@/snap userd
EnvironmentFile=-@SNAPD_ENVIRONMENT_FILE@
Restart=always

[Install]
WantedBy=default.target
