# data/systemd/snapd.unshare.service.in
[Unit]
Description=Apply rshared mount to @SNAP_MOUNT_DIR@ directory
RequiresMountsFor=/snap

[Service]
Type=oneshot
ExecStart=/bin/mkdir -p @SNAP_MOUNT_DIR@
ExecStart=/bin/mount -o bind @SNAP_MOUNT_DIR@ @SNAP_MOUNT_DIR@
ExecStart=/bin/mount --make-shared @SNAP_MOUNT_DIR@
ExecStop=/bin/umount -l @SNAP_MOUNT_DIR@
ExecStop=/bin/umount -l @SNAP_MOUNT_DIR@
RemainAfterExit=yes

[Install]
WantedBy=local-fs-pre.target
