# data/systemd/systemd-journald@snapd-security.service.d/00-snapd.conf
# Drop-in for the snapd-security journal namespace instance.
# Clears the default LogsDirectory to prevent failures with
# namespaced journald instances.
[Service]
LogsDirectory=
