# data/upstart/snap-userd.conf.in
start on started dbus

respawn
respawn limit 10 5

exec @bindir@/snap userd
