# data/upstart/snapd-user.conf.in
start on started dbus

respawn
respawn limit 10 5

exec @libexecdir@/snapd/snapd --user
