# core-initrd/latest/tee/usr/lib/systemd/system/run-mnt-tee\x2ddata.mount
[Unit]
Before=tee-supplicant.service

[Mount]
What=/dev/disk/by-partlabel/tee-data
Where=/run/mnt/tee-data
