# core-initrd/26.10/factory/usr/lib/systemd/system/initrd-udevadm-cleanup-db.service.d/override.conf
# https://github.com/systemd/systemd/pull/14341/files
[Unit]
Conflicts=systemd-udev-trigger.service