# core-initrd/26.04/factory/usr/lib/systemd/system/debug-shell.service.d/core-override.conf
[Unit]
AssertKernelCommandLine=dangerous

[Service]
ExecStart=
ExecStart=/bin/sh
