# debian/snapd.prerm
#! /bin/sh

set -e

if [ -d /run/systemd/system ]; then
    systemctl stop snapd.refresh.timer snapd.refresh.service snapd.autoimport.service snapd.socket snapd.service
    systemctl disable snapd.refresh.timer snapd.refresh.service snapd.autoimport.service snapd.socket snapd.service
fi
