# debian/snappy-autopilot.service
[Unit]
Description=Ubuntu Core Snappy Autopilot
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/snappy-go update
