# debian/snapd.frameworks.target
# all snappy framework services start before this; all snappy apps
# services start after it.

[Unit]
Description=Frameworks ready
