# debian/snappy-autopilot.service
[Unit]
Description=Ubuntu Core Snappy AutoUpdate
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/snappy update --automatic-reboot
