# debian/snapd.frameworks-pre.target
# all snappy framework services start after this

[Unit]
Description=Ready for frameworks
After=cloud-init.service

