# debian/snapd.frameworks-pre.target
# all snappy framework services start after this

[Unit]
Description=Ready for frameworks
After=snapd.firstboot.service
After=cloud-init.service

