# debian/ubuntu-snappy.run-hooks.service
[Unit]
Description=Run snappy compatibility hooks

[Service]
Type=oneshot
ExecStart=/usr/bin/snappy internal-run-hooks

[Install]
WantedBy=multi-user.target
