# debian/ubuntu-snappy.snapd.service
[Unit]
Description=Snappy daemon

[Service]
ExecStart=/usr/bin/snapd
