# debian/ubuntu-snappy.run-hooks.service
[Unit]
Description=Run snappy compatibility hooks
After=local-fs.target
Before=ubuntu-snappy.firstboot.service
DefaultDependencies=false

[Service]
Type=oneshot
ExecStart=/usr/bin/snappy internal-run-hooks

[Install]
WantedBy=multi-user.target
