# debian/ubuntu-snappy.snapd.service
[Unit]
Description=Snappy daemon
After=snappy-wait4network.service

[Service]
ExecStart=/usr/bin/snapd
