# debian/ubuntu-snappy.snapd.socket
[Unit]
Description=Socket activation for snappy daemon

[Socket]
ListenStream=@snapd

[Install]
WantedBy=sockets.target
