# etc/profile.d/apps-bin-path.sh
# Expand the $PATH to include /apps/bin which is what snappy applications
# use
PATH=$PATH:/apps/bin
