# data/systemd/snapd.service.d/security-journal.conf
# Pull in the snapd-security journal namespace socket alongside snapd.
# This is a best-effort dependency: if the socket is masked or fails,
# snapd starts normally without security logging.
[Unit]
Wants=systemd-journald@snapd-security.socket
After=systemd-journald@snapd-security.socket
