# core-initrd/26.04/fips/usr/lib/systemd/system/snap-initramfs-mounts.service.d/fips.conf
[Unit]
After=systemd-modules-load.service

[Service]
# Ensure tcrypt was loaded
ExecStartPre=modprobe tcrypt
# Then allow snap-bootstrap to mount things
# Self-test verify FIPS computed hmac
ExecStartPost=kcapi-hasher -n sha512hmac --check /run/mnt/kernel/.kernel.efi.hmac
# On failure, emergency.target is triggered with 10 second delay
