# core-initrd/26.04/factory/usr/lib/systemd/system/snapd.recovery-chooser-trigger.service.d/survive-switch-root.conf
[Unit]
# Because it is not killed by the switch root, we have to make sure
# not to not depend on local-fs.target, sysinit.target,
# etc. so that `isolate initrd-switch-root.target` will still take
# down all those targets.
# We still need to run after basic.target but this will be already wanted
# by initrd.target
DefaultDependencies=no
After=sysinit.target
After=basic.target

[Service]
ExecStart=
ExecStart=@/usr/lib/snapd/snap-bootstrap @snap-bootstrap recovery-chooser-trigger
